# home/start/disable-boost.service
[Unit]
Description=Disable Boost

[Service]
ExecStart=/home/sudanmerinosu/start/disable-boost.sh

[Install]
WantedBy=multi-user.target